# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=/bin/bash -c "cd /app/backend && node src/server.js"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.log
priority=1
environment=HOME="/root",NODE_ENV="production"

[program:frontend]
command=/bin/bash -c "sleep 5 && cd /app/frontend && node server.js"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend.log
priority=2
environment=HOME="/root",NODE_ENV="production",PORT="3000"
